# ansible/roles/edge/files/edge.service
[Unit]
Description=emitio edge
After=network.target

[Service]
ExecStart=/usr/local/bin/edge
Type=simple

[Install]
WantedBy=default.target